# pmx-config-example/pmx-config-configurator/onboot.sv.conf
[program:initialconfig]
command=/scripts/redirect-proxy-to-self.sh
autostart=true
autorestart=false
