# atlassian-agile/db/onboot.sv.conf
[program:onboot]
command=/scripts/onboot.sh
autostart=true
autorestart=false
